# linux_services/omv-power-alerts-wake.timer
[Unit]
Description=Wake detect timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
