# linux_services/omv-power-alerts-wake.service
[Unit]
Description=Notify on system wake

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/omv-power-alerts.d/wake_detect.sh
